# conf/supervisor.ini
[supervisord]
nodaemon = true
loglevel = DEBUG

[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
autorestart = true

[program:gunicorn]
command = gunicorn -w 2 -b 127.0.0.1:5000 -k gevent prediction_app.app
autorestart=true
